# supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=node bootstrap.js

[program:nginx]
command=nginx
